# launchd/com.tmact.statusd.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tmact.statusd</string>

  <key>ProgramArguments</key>
  <array>
    <string>__TMACT_BIN__</string>
    <string>statusd</string>
    <string>start</string>
    <string>--interval</string>
    <string>3s</string>
    <string>--state-path</string>
    <string>/tmp/tmact-status.json</string>
    <string>--tmux-options</string>
    <string>--log-path</string>
    <string>/tmp/tmact-statusd.jsonl</string>
    <string>--web-addr</string>
    <string>127.0.0.1:7890</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__TMACT_WORKDIR__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>__TMACT_PATH__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/tmact-statusd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/tmact-statusd.err.log</string>
</dict>
</plist>
